# crontab
SHELL=/bin/bash

*/15 * * * * root /opt/BGP-Table/bot_cronjob.sh &>> /opt/BGP-Table/bot_cronjob.log
